# argocd/applications/actualbudget/crontab
50 * * * * rsync -a /data/source/* /data/destination/
