# argocd/applications/actualbudget/crontab
0 */6 * * * rsync -a /data/source/* /data/destination/
